# deploy/bingx-scanner.service
[Unit]
Description=BingX TradFi realtime scanner
Documentation=https://github.com/USER/bingx-scanner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gapscan
Group=gapscan
WorkingDirectory=/opt/bingx-scanner
# Слушаем ТОЛЬКО loopback — наружу отдаёт Caddy (reverse-proxy + HTTPS + Basic Auth).
# Для Tailscale-режима bootstrap.sh заменит на GAP_HOST=0.0.0.0 (доступ режет ufw).
Environment=GAP_HOST=127.0.0.1
Environment=GAP_PORT=8787
Environment=GAP_OPEN_BROWSER=0
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
# Торговых ключей BingX тут НЕТ. FINNHUB_API_KEY (read-only котировки, close pc) — из отдельного
# файла /etc/bingx-scanner.env (права 600, в git НЕ коммитится). '-' = не падать, если файла нет.
EnvironmentFile=-/etc/bingx-scanner.env
ExecStart=/usr/bin/python3 /opt/bingx-scanner/bingx_gap_scanner.py
Restart=always
RestartSec=3
TimeoutStopSec=10

# --- hardening (пишем только в data/ для sqlite-лога сигналов) ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/bingx-scanner/data
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target
